# systemd/claudex.service
[Unit]
Description=Claudex - Claude Code Autonomous Agent (Telegram)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.claude-agent
ExecStart=%h/.claude-agent/scripts/start-claudex.sh
Restart=on-failure
RestartSec=30
Environment=HOME=%h
Environment=PATH=%h/.bun/bin:%h/.local/bin:%h/.cargo/bin:/usr/local/bin:/usr/bin:/bin
Environment=LANG=en_US.UTF-8
StandardOutput=append:%h/.claude-agent/logs/claudex-systemd.log
StandardError=append:%h/.claude-agent/logs/claudex-systemd.log

# Keep running after user logs out
KillMode=process
TimeoutStopSec=30

[Install]
WantedBy=default.target
